# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program:tail]
command=tail -f /app/log/script.log